# package/utils/busybox/files/root.crontab
0 * * * * /usr/sbin/ntpd -q -p kr.pool.ntp.org
